# container-binalerter.timer
[Unit]
Description=Daily check for putting bins out

[Timer]
OnCalendar=*-*-* 17:00:00
Persistent=true

[Install]
WantedBy=timers.target
